# backend/system_configuration/bakery.service
#!/bin/bash

# 1. Configuration
APP_ROOT="/home/nathanlooenping/bakery-platform"
BACKEND_DIR="$APP_ROOT/backend"
VENV_PYTHON="$APP_ROOT/venv/bin/python"
LOG_FILE="/var/log/webserver/logfile.log"
SSL_KEY="/home/nathanlooenping/cert.key"
SSL_CERT="/home/nathanlooenping/cert.pem"

# Create the log directory and file, then set permissions
mkdir -p /var/log/webserver
touch "$LOG_FILE"
chmod 666 "$LOG_FILE"

# 2. Create the systemd service file
cat <<EOF > /etc/systemd/system/bakery.service
[Unit]
Description=Uvicorn instance to serve Bakery Platform
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=$BACKEND_DIR
Environment="PYTHONUNBUFFERED=1"
Environment="KEY_PATH=$BACKEND_DIR/bucket-key.json"
Environment="BUCKET_NAME=cs5224-image-bucket"
Environment="PATH=$APP_ROOT/venv/bin"
ExecStart=$VENV_PYTHON -m uvicorn main:app --host 0.0.0.0 --port 443 --access-log

# Log both stdout and stderr to the same single file
StandardOutput=append:$LOG_FILE
StandardError=append:$LOG_FILE

Restart=always

[Install]
WantedBy=multi-user.target
EOF

# 3. Reload and Start
systemctl daemon-reload
systemctl enable bakery.service
systemctl restart bakery.service
